# scripts/launchd/com.photo.sync-albums.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.photo.sync-albums</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd /Users/jw-home/Work/photo_system/ai-photo-system &amp;&amp; PYTHONPATH=. /Users/jw-home/.local/bin/poetry run python scripts/sync_photos_albums.py --limit 1500 &gt;&gt; scripts/_inventory/sync_photos_albums.log 2&gt;&amp;1</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>35</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>/Users/jw-home/Work/photo_system/ai-photo-system</string>

    <key>StandardOutPath</key>
    <string>/Users/jw-home/Work/photo_system/ai-photo-system/scripts/_inventory/sync_photos_albums.stdout</string>

    <key>StandardErrorPath</key>
    <string>/Users/jw-home/Work/photo_system/ai-photo-system/scripts/_inventory/sync_photos_albums.stderr</string>
</dict>
</plist>
